# root/etc/services.d/boinc-client/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid root /usr/bin/boinc \
    --dir /data
